# templates/default/supervisord.fabric.service.include.conf.erb
[program:fabric_service]
command = /usr/bin/python /var/fabric_add_servers_monitor.py
autostart=true
autorestart=true

stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/fabric_service.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/fabric_service.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1

[group:fabric_service_server]
programs=fabric_service

