# templates/default/supervisord.fabric.conf.erb
[program:fabric]
command = /usr/bin/mysqlfabric manage start
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/fabric.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/fabric.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1


[group:fabric_server]
programs=fabric
